# dist/wesher.service
[Unit]
Description=wesher - wireguard mesh builder
After=network-online.target

[Service]
EnvironmentFile=-/etc/default/wesher
ExecStart=/usr/loca/sbin/wesher
Restart=on-failure
Type=simple

[Install]
WantedBy = multi-user.target